# i3/.config/i3/i3-wallpaper.cron
# m h dom mon dow command
@hourly ~/.i3/bing_wallpaper.sh en-IN true
#@hourly wallpaper-reddit -f
